# packages/traefik/traefik.service
[Unit]
Description={{.Description}}
Documentation={{.URL}}
After=network-online.target
Wants=network-online.target

[Service]
User={{.Name}}
Group={{.Name}}
ExecStart=/usr/bin/{{.Name}} /etc/{{.Name}}/{{.Name}}.toml
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
Restart=on-failure
