# packages/docker-gc/docker-gc-volumes.timer
[Unit]
Description=Timer for docker-gc-volumes
After=docker-gc-images.service

[Timer]
Unit=docker-gc-volumes.service
OnBootSec=10min
OnUnitActiveSec=8h

[Install]
WantedBy=multi-user.target
